# infra/seed/botho-seed.service
[Unit]
Description=Botho Seed Node (Testnet)
Documentation=https://github.com/botho-project/botho
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu

# Run as testnet relay node (no minting on seed node).
#
# This is the correct unit for a MULTI-SEED network (>= 2 seeds): seeds relay
# and gossip, while minting happens on a separate validator/faucet node, and
# the [network.quorum] mode stays "recommended" (auto BFT threshold).
#
# SINGLE-SEED bootstrap exception: a lone seed cannot form a "recommended"
# quorum (it needs >= 2 nodes -> minting stalls at "have 1, need 2"). To mint
# from a single self-quorum seed during initial bring-up, instead use:
#
#   ExecStart=/usr/local/bin/botho --testnet run --mint --mint-threads 1
#
# and set the following in ~/.botho/testnet/config.toml so the lone node forms
# its own quorum:
#
#   [network.quorum]
#   mode = "explicit"
#   threshold = 1
#   members = []
#
# Switch BACK to relay mode (no --mint flags, [minting] enabled = false in
# config.toml, quorum mode "recommended") once a second seed / validator
# joins. See infra/seed/README.md "Single-seed vs multi-seed".
#
# NOTE: there is no --relay CLI flag (verified against the 3.0.0 binary,
# 2026-07-04 — passing it fails with "unexpected argument"). Relay behavior
# is simply running without minting.
ExecStart=/usr/local/bin/botho --testnet run

# Restart on failure
Restart=on-failure
RestartSec=10

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/.botho

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=botho-seed

[Install]
WantedBy=multi-user.target
